# bin/launchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aptgregator.tek.production</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/tek/src/aptgregator/bin/initd.sh</string>
    <string>run</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>UserName</key>
  <string>tek</string>
  <key>WorkingDirectory</key>
  <string>/Users/tek/src/aptgregator</string>
  <key>StandardErrorPath</key>
  <string>/Users/tek/src/aptgregator/var/log/launchAgent.log</string>
  <key>StandardOutPath</key>
  <string>/Users/tek/src/aptgregator/var/log/launchAgent.log</string>
</dict>
</plist>

